# files/maverick/etc/update-motd.d/00-header
#!/bin/sh

#uname -a
#printf "%s\n" "$(lsb_release -s -d)"

uname -a | /usr/games/cowthink
